# config/supervisord/supervisord.ini
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n
autostart=true
autorestart=true

[program:cron]
command = /usr/sbin/crond -f -L 15
autostart=true
autorestart=true

[program:logstash]
command=/opt/logstash/bin/logstash agent --config /etc/logstash/conf.d/ --log /var/log/logstash/logstash.log
autostart=true
autorestart=true
user=logstash
